# popays-bot.service
[Unit]
Description=POPAYS Fast Food Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/popays-bot
Environment=PATH=/opt/popays-bot/venv/bin
ExecStart=/opt/popays-bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Logging
SyslogIdentifier=popays-bot

[Install]
WantedBy=multi-user.target
